# stuff/systemd/pingstat.service
[Unit]
Description=Pingstat daemon

[Service]
Type=simple
ExecStart=/opt/pingstat/pingstat

[Install]
WantedBy=multi-user.target
